# ops/launchd/com.agent44.smoke-nyk-nav.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Hourly NYK arrow-nav regression check on the Mac mini runner.

  - Fires every hour at :00, but the command sleeps a uniformly random
    0–1799 seconds before triggering the workflow. Final fire times are
    spread randomly across each hour so the cadence doesn't look like a
    bot signing in on the dot.
  - Runs headed (HEADFUL=true) for a less detectable Chromium fingerprint.
    Requires the agent44 user to be logged in on the Mac mini (display
    can sleep — only the user session needs to be active).
  - Test failures email botwhisperer@hey.com only (NYK_SMOKE_RECIPIENTS in
    the workflow's GH secret should be set to that address). Lora is never
    notified by this test.

  Install:
    cp ops/launchd/com.agent44.smoke-nyk-nav.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.agent44.smoke-nyk-nav.plist

  After verifying it fires correctly, retire the legacy plist:
    launchctl unload ~/Library/LaunchAgents/com.agent44.smoke-nyk.plist
    rm ~/Library/LaunchAgents/com.agent44.smoke-nyk.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agent44.smoke-nyk-nav</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>sleep $((RANDOM % 1800)); /opt/homebrew/bin/gh workflow run smoke-nyk.yml --ref main -R agent44bot/agent44_app -f test=nav -f headful=true</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/agent44</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>

    <!-- Fire every hour at :00. Command-level sleep adds 0–30 min jitter. -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Minute</key><integer>0</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/tmp/smoke-nyk-nav.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/smoke-nyk-nav.err.log</string>
</dict>
</plist>
